# pdf_server/deploy/mymserver-tunnel.service
[Unit]
Description=MyMultiplatform Public Tunnel (Cloudflare Quick Tunnel)
After=network-online.target nginx.service mymserver.service
Wants=network-online.target

[Service]
Type=simple
User=des
Group=des
WorkingDirectory=/home/des
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:80 --no-autoupdate
Restart=always
RestartSec=5
StandardOutput=append:/home/des/mymserver-cloudflared.log
StandardError=append:/home/des/mymserver-cloudflared.log

[Install]
WantedBy=multi-user.target
